# mpuserver.service
[Service]
Type=simple
ExecStart=/bin/mpuserver
User=mpu

[Install]
WantedBy=multi-user.target
